# konzentrator/rc.local
#!/bin/sh -e
# rc.local
/opt/eulenfunk/konzentrator/bgp-konzentrator-setup.sh
exit 0

